# systemd/climate-update.service
[Unit]
Description=Refresh downloaded climate images
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/climate pull --update
